# etc/systemd/system/mnt-data.mount
[Mount]
What=/dev/disk/by-label/data
Where=/mnt/data
Type=ext4
Options=noatime
